# scripts/sing-box-release.service
[Unit]
Description=Authorization of sing-box Distribution
After=network.target

[Service]
User=shuqixiao
Group=shuqixiao
WorkingDirectory=/home/shuqixiao/Projects/sing-box-docker
ExecStart=/home/shuqixiao/.local/bin/uv run /home/shuqixiao/Projects/sing-box-docker/auth_service.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
